# system/etc/systemd/system/greetd.service.d/10-dotfiles-vtrgb.conf
[Service]
ExecStartPre=/usr/bin/setvtrgb /etc/vtrgb-dotfiles
